# env/dev/supervisor.conf
[program:api]
command=./bin/mix api > /dev/null 2>&1 &
numprocs=1
startretries=3
